# projects/Amlogic/filesystem/usr/lib/udev/rules.d/80-bluetooth-uart.rules
KERNEL=="ttyS1", SUBSYSTEMS=="tty", ENV{SYSTEMD_WANTS}="bluetooth.service"


KERNEL=="ttyS1", SUBSYSTEMS=="tty", RUN+="modprobe hci-uart"